# notification-client/systemd/notification-client.service
[Unit]
Description=Notification Client Service
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=nobody
Group=nogroup
ExecStart=/usr/bin/notification-client -d -c /etc/notification-client/config.json
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=notification-client

[Install]
WantedBy=multi-user.target