# debian/player-client.service
[Unit]
Description=PlayerClient - HTTP Audio Playback Server
Documentation=https://github.com/playerclient/player-client
After=network.target sound.target

[Service]
Type=simple
ExecStart=/usr/local/bin/player-client
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
