# deb/DEBIAN/postinst
#!/bin/bash
# postinst script for notification-player

set -e

echo "Notification Player post-installation configuration..."

# 创建用户
if ! id -u notification-player >/dev/null 2>&1; then
    useradd -r -s /bin/false -d /var/lib/notification-player notification-player
    echo "Created notification-player user"
fi

# 创建目录
mkdir -p /var/lib/notification-player
mkdir -p /var/log/notification-player
mkdir -p /usr/share/notification-player/sounds

# 设置权限
chown -R notification-player:notification-player /var/lib/notification-player
chown -R notification-player:notification-player /var/log/notification-player
chown -R notification-player:notification-player /usr/share/notification-player

# 生成默认音频文件（如果有工具）
if command -v sox >/dev/null 2>&1; then
    sox -n /usr/share/notification-player/sounds/default.wav synth 0.5 sine 800
elif command -v ffmpeg >/dev/null 2>&1; then
    ffmpeg -f lavfi -i "sine=frequency=800:duration=0.5" -ar 44100 -ac 2 \
           /usr/share/notification-player/sounds/default.wav -y 2>/dev/null
fi

# 启用服务
systemctl daemon-reload
systemctl enable notification-player.service
systemctl start notification-player.service || true

echo "Post-installation complete"
exit 0